# packages/scheduler/cronTabFile/crontab
0 0 * * * /usr/scheduler/jobs/runVSExport.sh

0 0 * * * /usr/scheduler/jobs/refreshMetabaseData.sh